# libs/coder-devcontainer/supervisord.conf
[supervisord]
nodaemon=false
logfile=/workspaces/setup/logs/supervisor/supervisord.log
pidfile=/workspaces/setup/logs/supervisor/supervisord.pid
user=coder

[program:fleet-mcp]
directory=/workspaces/setup
command=/workspaces/setup/bin/direnv exec . nx run --tui=false fleet-mcp:server
user=coder
autostart=true
autorestart=true
startretries=3
stopasgroup=true
killasgroup=true
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/workspaces/setup/logs/fleet-mcp/server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5

[inet_http_server]
port=127.0.0.1:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
