# libs/coder-devcontainer/supervisord.conf
[supervisord]
nodaemon=false
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
user=root

[program:fleet-mcp]
directory=/workspaces/setup
command=direnv exec . nx run --tui=false fleet-mcp:server
user=coder
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/fleet-mcp/server.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
environment=HOME="/home/coder",USER="coder"

[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface
